# packages/linux-clipboard-manager/linux_clipboard_manager-1.0.7.tar.gz/linux_clipboard_manager-1.0.7/smart-clipboard.service
[Unit]
Description=Smart Clipboard Manager
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/smart-clipboard
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
