# packages/sollol/sollol-0.9.66.tar.gz/sollol-0.9.66/systemd/sollol-rpc-server.service
[Unit]
Description=SOLLOL llama.cpp RPC Server
After=network.target
Documentation=https://github.com/BenevolentJoker-JohnL/SynapticLlamas

[Service]
Type=simple
User=%u
WorkingDirectory=%h
ExecStart=%h/.local/bin/rpc-server --host 0.0.0.0 --port 50052
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
Nice=-5

[Install]
WantedBy=default.target
