# packages/gorunn/gorunn-0.1.10-py3-none-any.whl/gorunn/templates/mounts/node/files/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/dev/null
logfile_maxbytes=0
loglevel = INFO

[program:node]
command=/scripts/node.sh
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
